# dqxtool-api.service
[Unit]
Description=DQX Blacksmith Simulator API Server (ONNX Inference Backend)
After=network.target

[Service]
Type=simple
# 【重要】本番サーバー上に存在する正しいユーザー名に書き換えてください。
# (現在操作しているユーザー名は Linux 上で `whoami` コマンドを実行することで確認できます)
User=ubuntu

# 【重要】Gitリポジトリ（DQXTool）が配置されている実際の絶対パスを指定してください
WorkingDirectory=/home/ubuntu/DQXTool

# npm の絶対パスを指定してください (which npm で確認できます)
ExecStart=/usr/bin/npm run server

# サーバーが落ちた場合に自動で再起動する設定
Restart=on-failure
RestartSec=5

# 環境変数の設定 (ポート番号などを指定可能)
Environment=PORT=3001 NODE_ENV=production

# ログをシステムジャーナルに出力
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
